# Users/rajes/OneDrive/Desktop/Auto post/files (12)/autopost_bot.service
[Unit]
Description=Telegram AutoPosterBot
After=network.target mongod.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/telegram_autopost
ExecStart=/opt/telegram_autopost/venv/bin/python -m autobot.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
# Pass the .env file so the service reads env vars
EnvironmentFile=/opt/telegram_autopost/.env

[Install]
WantedBy=multi-user.target
